# build/receiptd.rc.d
#!/bin/ksh

daemon="/usr/local/sbin/httpftd" 
. /etc/rc.d/rc.subr

rc_cmd $1